# scripts/systemd/gerbera-mysql.service
[Unit]
Description=Gerbera Media Server with MySQL
After=mysql.target network.target

[Service]
Type=simple
User=gerbera
Group=gerbera
ExecStart=/usr/bin/gerbera -c /etc/gerbera/config.xml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
